# funipro-api.service
[Unit]
Description=FuniPro API .NET Web API App running on Ubuntu
After=network.target mssql-server.service

[Service]
Type=notify
# Nome do usuário que executará o serviço
User=www-data
# Diretório de trabalho (onde está o executável publicado)
WorkingDirectory=/var/www/funipro/backend/publish
# Comando para iniciar o serviço
ExecStart=/usr/local/bin/dotnet /var/www/funipro/backend/publish/FuniproApi.dll
# Reiniciar serviço em caso de falha
Restart=always
# Aguardar 10 segundos antes de reiniciar
RestartSec=10
# Sinal de encerramento
KillSignal=SIGINT
# Identificador para logs
SyslogIdentifier=funipro-api
# Variáveis de ambiente
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://localhost:5000

[Install]
WantedBy=multi-user.target

